# log_parser_install/systemd/log_parser.service
[Unit]
Description=Script for parsing nginx logs and cleanup them

[Service]
ExecStart=/opt/log_parser.sh
Restart=always
RestartSec=3